# .cron
0 */12 * * * mu index
0 */4 * * * vdirsyncer sync
